# restic-prune.target
[Unit]
Description=Completed restic prune tasks
StopWhenUnneeded=yes
After=restic-forget.target
