# restic-forget@.service.in
[Unit]
Description=restic forget for %i
After=restic-backup@%i.service

[Service]
Type=oneshot
User=@RESTIC_USER@
Environment=RESTIC_LOCKFILE=/run/restic/forget.lock
Environment=XDG_CACHE_HOME=/var/cache/restic
EnvironmentFile=-/etc/restic/restic.conf
EnvironmentFile=/etc/restic/%i/restic.conf
WorkingDirectory=/etc/restic/%i
IOSchedulingClass=idle
ExecStart=/usr/bin/flock --verbose $RESTIC_LOCKFILE \
	@RESTIC_PATH@ $RESTIC_COMMON_ARGS \
	forget $RESTIC_FORGET_ARGS

[Install]
WantedBy=multi-user.target
