# restic-forget.target
[Unit]
Description=Completed restic forget tasks
StopWhenUnneeded=yes
After=restic-backup.target
