# restic-prune-schedule.timer
[Unit]
Description=@schedule@ prune of %i

[Timer]
OnCalendar=@schedule@
Unit=restic-prune@%i.service

[Install]
WantedBy=timers.target
