# restic-prune.target
[Unit]
Description=All prune tasks should be Before=restic-forget.target
After=restic-forget.target
