# restic-clean-schedule.timer
[Unit]
Description=clean @schedule@ backups of %I

[Timer]
OnCalendar=@schedule@
Unit=restic-clean@%i.service

[Install]
WantedBy=timers.target
