# restic-backup.target
[Unit]
Description=Completed restic backup tasks
StopWhenUnneeded=yes
