# restic-forget.target
[Unit]
Description=All forget tasks should be Before=restic-forget.target
After=restic-backup.target
