# restic-forget-schedule.timer
[Unit]
Description=@schedule@ forget of %i

[Timer]
OnCalendar=@schedule@
Unit=restic-forget@%i.service

[Install]
WantedBy=timers.target
